# rootfs/etc/s6-overlay/s6-rc.d/uwsgi-cgit/run
#!/command/execlineb -P

/usr/sbin/uwsgi \
  --uid git --gid git \
  --master --die-on-term --listen 256 \
  --processes 1 --threads 2 \
  --logto /config/logs/uwsgi/cgit.log \
  --plugin cgi --cgi /app/cgit.cgi \
  --procname-master "uwsgi cgit" \
  --hook-master-start "unix_signal:1 gracefully_kill_them_all" \
  --socket /run/uwsgi/cgit.socket
