# rootfs/etc/s6-overlay/s6-rc.d/openssh/run
#!/bin/sh
exec 2>&1
exec /usr/sbin/sshd -D -e -f /config/etc/ssh/sshd_config
