# rootfs/etc/s6-overlay/s6-rc.d/openssh/run
#!/command/execlineb -P

/usr/sbin/sshd -D -e -f /config/etc/ssh/sshd_config
