# rootfs/etc/s6-overlay/s6-rc.d/uwsgi-cgit-log/run
#!/bin/sh
exec logutil-service /config/log/uwsgi-cgit
