# rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/command/execlineb -P

/usr/sbin/crond -f
