# rootfs/etc/s6-overlay/s6-rc.d/uwsgi-git-http/run
#!/command/execlineb -P

/usr/sbin/uwsgi \
  --uid git --gid git \
  --master --die-on-term --listen 256 \
  --processes 1 --threads 2 --post-buffering 8192 \
  --env GIT_PROJECT_ROOT=/config/repositories --env HOME=/config/repositories \
  --logto /config/logs/uwsgi/git-http-backend.log \
  --plugin cgi --cgi /usr/libexec/git-core/git-http-backend \
  --procname-master "uwsgi git-http-backend" \
  --hook-master-start "unix_signal:1 gracefully_kill_them_all" \
  --socket /run/uwsgi/git-http-backend.socket
